# mybookmarks-go.cron
#! /bin/bash

/var/www/mybookmarks-go/mybookmarks-go backup