# scripts/mybookmarks.cron
#! /bin/bash

/var/www/mybookmarks/mybookmarks backup
